# scripts/init/opx-ip.service
[Unit]
Description=IP address handler
After=opx-cps.service
Wants=opx-cps.service
OnFailure=service_onfailure@%n.service

[Service]
Type=notify
EnvironmentFile=/etc/opx/opx-environment
ExecStart=/usr/bin/python -u /usr/bin/base_ip.py

[Install]
WantedBy=multi-user.target
